# root/etc/services.d/dvblink/run
#!/usr/bin/with-contenv bash

LD_LIBRARY_PATH=/usr/local/bin/dvblink/lib /usr/local/bin/dvblink/dvblink_server -command_line_mode
